# .config/systemd/user/gcalcli-remind.timer
[Unit]
Description=Run gcalcli reminders every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=1min

[Install]
WantedBy=timers.target
