# .config/systemd/user/gcalcli-remind.service
[Unit]
Description=Google Calendar reminders via gcalcli
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/home/linn/.local/bin/experimental/gcal-remind
